# profile.d/env.sh
PATH=$PATH:$HOME/bin
PS1="$(whoami)=; "

export PATH PS1
