# profile.d/git.sh
alias gp='git push origin $(git rev-parse --abbrev-ref HEAD)'
alias gcd='cd $(git rev-parse --show-topleve)'
